# ops/systemd/drone-v2-main-deploy.timer
[Unit]
Description=Poll origin/main and deploy drone_v2 when it changes

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=15s
Persistent=true

[Install]
WantedBy=timers.target
